# config/com.influxdb.optimized.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.influxdb.optimized</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/influxd</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>INFLUXD_STORAGE_CACHE_MAX_MEMORY_SIZE</key>
        <string>2147483648</string>
        <key>INFLUXD_STORAGE_CACHE_SNAPSHOT_MEMORY_SIZE</key>
        <string>52428800</string>
        <key>INFLUXD_STORAGE_MAX_CONCURRENT_COMPACTIONS</key>
        <string>8</string>
        <key>INFLUXD_STORAGE_SERIES_ID_SET_CACHE_SIZE</key>
        <string>104857600</string>
        <key>INFLUXD_STORAGE_WAL_FSYNC_DELAY</key>
        <string>100ms</string>
        <key>INFLUXD_STORAGE_WAL_MAX_CONCURRENT_WRITES</key>
        <string>100</string>
        <key>INFLUXD_LOG_LEVEL</key>
        <string>info</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/opt/homebrew/var/log/influxdb-optimized.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/homebrew/var/log/influxdb-optimized.error.log</string>
</dict>
</plist>